# private_dot_config/systemd/user/ulauncher.service.d/override.conf
[Service]
Environment=GDK_BACKEND=wayland
ExecStart=
ExecStart=/usr/bin/ulauncher --hide-window
