# private_dot_config/systemd/user/random-wallpaper.service
[Unit]
Description=Set the desktop background

[Service]
Type=oneshot
ExecStart=/home/dan/.local/bin/random-wallpaper

